# plantgateway_conf/supervisor_conf/rtl_433.conf
[program:rtl_433_daemon]
command=rtl_433 -F "mqtt://192.168.1.19:1883" -Y autolevel -M level -M noise
directory=/home/pi
autostart=true
autorestart=true
user=pi

[program:rtl_433_hass_mqtt_discovery_daemon]
command=/home/pi/devel/pyenvs/rtl_433_mqtt_hass/bin/python /home/pi/rtl_433/examples/rtl_433_mqtt_hass.py --retain --host 192.168.1.19 --ignore "tpms,interlogix-security"
directory=/home/pi
autostart=true
autorestart=true
user=pi 
