# plantgateway_conf/supervisor_conf/x10_heyu_mqtt.conf
[program:x10_heyu_mqtt]
command=sudo -u pi python3 /home/pi/x10_heyu_mqtt/run.py
directory=/home/pi/x10_heyu_mqtt
autostart=true
autorestart=true
user=root
